# arch-sshd/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s
numprocs=1
directory=/
autostart=true
autorestart=true
startsecs=10
startretries=3
user=root
stdout_logfile=/var/log/supervisor/sshd-stdout.log
stderr_logfile=/var/log/supervisor/sshd-stderr.log
