# OLD/arch-nginx-php-mariadb-wordpress/wordpress.ini
[program:wordpress]
command=/opt/wordpress-start.sh
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=false
user=root

