# arch-nginx-php-mariadb-ttrss/ttrss-cron.ini
[program:ttrss]
command=/bin/bash -c 'while true; do php /srv/http/update.php --feeds --quiet; sleep 900; done'
process_name=%(program_name)s
numprocs=1
directory=/srv/http
autostart=true
autorestart=true
startsecs=10
startretries=3
user=http
stdout_logfile=/var/log/supervisor/ttrss-stdout.log
stderr_logfile=/var/log/supervisor/ttrss-stderr.log
