# OLD/arch-plex/plex.ini
[program:plex]
command=/opt/plexmediaserver/start_pms
process_name=%(program_name)s
numprocs=1
directory=/opt/plexmediaserver
autostart=true
autorestart=true
startsecs=10
startretries=3
user=root
stdout_logfile=/var/log/supervisor/plex-stdout.log
stderr_logfile=/var/log/supervisor/plex-stderr.log
